# deploy/vps/systemd/tg-msg-manager-update.timer
[Unit]
Description=Run tg-msg-manager update hourly

[Timer]
OnCalendar=*-*-* *:00:00 Europe/Kyiv
AccuracySec=1s
RandomizedDelaySec=0
Unit=tg-msg-manager-update.service

[Install]
WantedBy=timers.target
